# deploy/arxiv-agent.timer
[Unit]
Description=Run arxiv-agent daily at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
