# deploy/essay-agent.service
[Unit]
Description=essay-agent daily paper fetcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/essay-agent
ExecStart=/opt/essay-agent/.venv/bin/python essay_agent.py
Environment=PYTHONUNBUFFERED=1
MemoryMax=768M
MemoryHigh=512M
CPUQuota=80%
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
